# cronjobs/crontab
0,*/30 9-17 * * 1-5 ~/scripts/git-auto-fetch > /dev/null
